# assets/scripts/splunk_health_check.cron
*/10 * * * * root /opt/splunkforwarder/bin/splunk status || /opt/splunkforwarder/bin/splunk start
